# scripts/macos_direct/pkg/postinstall
#!/bin/zsh
# Copyright (c) EZBLOCK Inc & AUTHORS
# SPDX-License-Identifier: BSD-3-Clause

set -euo pipefail

PLIST_NAME="io.cylonix.sase.direct.daemon"
PLIST_DST="/Library/LaunchDaemons/$PLIST_NAME.plist"
NOTIFIER_PLIST_NAME="io.cylonix.sase.direct.notifier"
NOTIFIER_PLIST_DST="/Library/LaunchAgents/$NOTIFIER_PLIST_NAME.plist"
APP_PATH="/Applications/Cylonix.app"
DIRECT_BUNDLE_ID="io.cylonix.sase.direct"
DAEMON_PATH="$APP_PATH/Contents/Resources/cylonixd"
NOTIFIER_PATH="$APP_PATH/Contents/Resources/CylonixNotifier.app/Contents/MacOS/CylonixNotifier"

# ---------------------------------------------------------------------------
# Flavor/orphan hygiene (keep in sync across pkg/preinstall, pkg/postinstall
# and uninstall_direct.sh). See pkg/preinstall for the full story: a
# different-flavor bundle at $APP_PATH makes Installer dodge the destination
# (BundleHasStrictIdentifier) into "/Applications/Cylonix*.localized/", and
# the orphans' share extensions break the Share menu.
# ---------------------------------------------------------------------------
LSREGISTER="/System/Library/Frameworks/CoreServices.framework/Versions/A/Frameworks/LaunchServices.framework/Versions/A/Support/lsregister"

cylonix_bundle_id() {
  /usr/libexec/PlistBuddy -c 'Print :CFBundleIdentifier' \
    "$1/Contents/Info.plist" 2>/dev/null || true
}

remove_cylonix_bundle() {
  local app="$1" gui_user
  gui_user="$(stat -f '%Su' /dev/console 2>/dev/null || true)"
  if [[ -n "${gui_user:-}" && "$gui_user" != "root" && "$gui_user" != "_"* \
        && -d "$app/Contents/PlugIns/ShareExtension.appex" ]]; then
    sudo -u "$gui_user" pluginkit -r \
      "$app/Contents/PlugIns/ShareExtension.appex" 2>/dev/null || true
  fi
  "$LSREGISTER" -u "$app" 2>/dev/null || true
  rm -rf "$app"
}

remove_cylonix_orphans() {
  local removed=1 dir app id
  for dir in /Applications/Cylonix*.localized(N); do
    app="$dir/Cylonix.app"
    id="$(cylonix_bundle_id "$app")"
    if [[ ! -d "$app" || -z "$id" || "$id" == io.cylonix.sase* ]]; then
      if [[ -d "$app" ]]; then
        remove_cylonix_bundle "$app"
      fi
      rm -rf "$dir"
      removed=0
    fi
  done
  return $removed
}

# Restart the GUI user's Finder so ShareKit drops its cached share-service
# identities. A plain `killall Finder` is silently ineffective from the
# PackageKit script sandbox (the signal never lands and stderr is
# discarded), so go through launchd — kickstart -k restarts the
# com.apple.Finder agent — keeping killall as a fallback for contexts
# without a GUI session.
bounce_finder() {
  local gui_user gui_uid
  gui_user="$(stat -f '%Su' /dev/console 2>/dev/null || true)"
  if [[ -n "${gui_user:-}" && "$gui_user" != "root" && "$gui_user" != "_"* ]]; then
    gui_uid="$(id -u "$gui_user" 2>/dev/null || true)"
    if [[ -n "${gui_uid:-}" ]] \
        && launchctl kickstart -k "gui/${gui_uid}/com.apple.Finder" 2>/dev/null; then
      return 0
    fi
  fi
  killall Finder 2>/dev/null || true
}

# Self-heal a strict-identifier dodge. preinstall clears the destination up
# front, so normally the payload is already at $APP_PATH. If Installer still
# diverted it, the FRESH install is sitting inside a Cylonix*.localized orphan
# while $APP_PATH holds the old different-flavor bundle — move the new bundle
# to the canonical path (newest orphan first) BEFORE the daemon below is
# pointed at $APP_PATH, then sweep the remaining litter.
if [[ "$(cylonix_bundle_id "$APP_PATH")" != "$DIRECT_BUNDLE_ID" ]]; then
  for dir in /Applications/Cylonix*.localized(Nom); do
    if [[ "$(cylonix_bundle_id "$dir/Cylonix.app")" == "$DIRECT_BUNDLE_ID" ]]; then
      echo "postinstall: payload was diverted to $dir; restoring to $APP_PATH" >&2
      if [[ -d "$APP_PATH" ]]; then
        killall Cylonix 2>/dev/null || true
        remove_cylonix_bundle "$APP_PATH"
      fi
      mv "$dir/Cylonix.app" "$APP_PATH"
      rm -rf "$dir"
      break
    fi
  done
fi
if remove_cylonix_orphans; then
  # ShareKit caches share-service identities inside Finder; bounce it so the
  # Share menu re-resolves against the surviving registration.
  bounce_finder
fi

# Create required directories
mkdir -p /var/lib/cylonix
mkdir -p /var/run/cylonix
mkdir -p /var/log/cylonix
chmod 755 /var/log/cylonix
# Ensure log file is world-readable so the app can display daemon logs
touch /var/log/cylonix/cylonixd.log
chmod 644 /var/log/cylonix/cylonixd.log

# Stop any existing daemon and WAIT for launchd to fully drop the label before we
# rewrite the plist. Deregister BY LABEL (bootout) rather than by plist path
# (`launchctl unload "$PLIST_DST"`): unload-by-path silently fails when a prior
# uninstall already removed the plist, leaving the label registered. bootout is
# asynchronous for a KeepAlive job, so poll until the label is actually gone —
# otherwise the bootstrap below hits EIO ("already loaded") and the daemon never
# starts (the uninstall→reinstall failure mode).
launchctl bootout "system/$PLIST_NAME" 2>/dev/null || true
i=0; while launchctl print "system/$PLIST_NAME" >/dev/null 2>&1 && [ "$i" -lt 100 ]; do
  i=$((i + 1)); sleep 0.1
done

# Install launchd plist
cat > "$PLIST_DST" <<'PLIST'
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>io.cylonix.sase.direct.daemon</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Applications/Cylonix.app/Contents/Resources/cylonixd</string>
		<string>--tun</string>
		<string>utun</string>
		<string>--state</string>
		<string>/var/lib/cylonix/cylonixd.state</string>
		<string>--statedir</string>
		<string>/var/lib/cylonix</string>
		<string>--socket</string>
		<string>/var/run/cylonix/cylonixd.sock</string>
		<string>--port</string>
		<string>41741</string>
	</array>
	<key>EnvironmentVariables</key>
	<dict>
		<key>CYLONIX_TAILDROP_DEFAULT_DOWNLOADS</key>
		<string>1</string>
	</dict>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>StandardOutPath</key>
	<string>/var/log/cylonix/cylonixd.log</string>
	<key>StandardErrorPath</key>
	<string>/var/log/cylonix/cylonixd.log</string>
</dict>
</plist>
PLIST

chown root:wheel "$PLIST_DST"
chmod 644 "$PLIST_DST"

# Load and start the daemon by LABEL. Do NOT suppress errors — a failure must be
# visible in /var/log/install.log, not a silently dead daemon. If the first
# bootstrap fails (e.g. a stubborn stale registration surviving the stop above),
# bootout + wait once more and retry before falling back to the legacy load.
if ! launchctl bootstrap system "$PLIST_DST"; then
  echo "postinstall: bootstrap failed; clearing stale/ghost registration and retrying" >&2
  # EIO with `print` reporting the label absent is a launchd "ghost": bootout has
  # nothing to grab but bootstrap collides with a cached entry. Evict it via every
  # available path — bootout by label, bootout by plist path, and the legacy
  # `remove` (a different code path that can drop a cached job bootout can't) —
  # then clear any disable override before retrying.
  launchctl bootout "system/$PLIST_NAME" 2>/dev/null || true
  launchctl bootout system "$PLIST_DST" 2>/dev/null || true
  launchctl remove "$PLIST_NAME" 2>/dev/null || true
  launchctl enable "system/$PLIST_NAME" 2>/dev/null || true
  i=0; while launchctl print "system/$PLIST_NAME" >/dev/null 2>&1 && [ "$i" -lt 100 ]; do
    i=$((i + 1)); sleep 0.1
  done
  launchctl bootstrap system "$PLIST_DST" || launchctl load "$PLIST_DST" || true
fi
launchctl kickstart -kp "system/$PLIST_NAME" 2>/dev/null || true

# Symlink CLI to /usr/local/bin so users can run "cylonix" from the terminal
mkdir -p /usr/local/bin
ln -sf "$APP_PATH/Contents/Resources/cylonix" /usr/local/bin/cylonix

# Note: the daemon resolves the active GUI/console user itself for
# Taildrop's Downloads/Cylonix destination (see
# tailscale/feature/taildrop/paths_cylonix.go), gated on
# CYLONIX_TAILDROP_DEFAULT_DOWNLOADS=1 (set above). No postinstall
# `cylonix set --operator` is required; admins on a shared Mac can still
# pin a specific user explicitly via `cylonix set --operator <user>`.

# Install the per-user notifier LaunchAgent. The system-wide LaunchAgents
# directory loads automatically for every user at login; we additionally
# bootstrap it into the currently-logged-in GUI session so notifications
# start working immediately without requiring logout/login.
cat > "$NOTIFIER_PLIST_DST" <<NOTIFIERPLIST
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>$NOTIFIER_PLIST_NAME</string>
	<key>ProgramArguments</key>
	<array>
		<string>$NOTIFIER_PATH</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>ThrottleInterval</key>
	<integer>5</integer>
	<key>StandardOutPath</key>
	<string>/tmp/cylonix-notifier.log</string>
	<key>StandardErrorPath</key>
	<string>/tmp/cylonix-notifier.log</string>
</dict>
</plist>
NOTIFIERPLIST
chown root:wheel "$NOTIFIER_PLIST_DST"
chmod 644 "$NOTIFIER_PLIST_DST"

GUI_USER="$(stat -f '%Su' /dev/console 2>/dev/null || true)"
if [[ -n "${GUI_USER:-}" && "$GUI_USER" != "root" && "$GUI_USER" != "_"* ]]; then
  GUI_UID="$(id -u "$GUI_USER" 2>/dev/null || true)"
  if [[ -n "${GUI_UID:-}" ]]; then
    # Replace any stale registration; ignore errors from never-loaded.
    launchctl bootout "gui/$GUI_UID/$NOTIFIER_PLIST_NAME" 2>/dev/null || true
    launchctl bootstrap "gui/$GUI_UID" "$NOTIFIER_PLIST_DST" 2>/dev/null || true
    launchctl enable "gui/$GUI_UID/$NOTIFIER_PLIST_NAME" 2>/dev/null || true
  fi
fi

# Kill any running Cylonix app so the new version is picked up
killall Cylonix 2>/dev/null || true
killall "cylonix-direct" 2>/dev/null || true

# Register app with LaunchServices so macOS uses this copy for share extension
LSREGISTER="/System/Library/Frameworks/CoreServices.framework/Versions/A/Frameworks/LaunchServices.framework/Versions/A/Support/lsregister"

# Unregister any stale copies from derived-data directories
"$LSREGISTER" -u "/tmp/cylonix-direct-derived" 2>/dev/null || true

# Register the installed app
"$LSREGISTER" -f "$APP_PATH"

# Register share extension with pluginkit
SHARE_EXT="$APP_PATH/Contents/PlugIns/ShareExtension.appex"
if [[ -d "$SHARE_EXT" ]]; then
  pluginkit -a "$SHARE_EXT" 2>/dev/null || true
fi

exit 0
