# scripts/macos_direct/pkg/postinstall
#!/bin/zsh
# Copyright (c) EZBLOCK Inc & AUTHORS
# SPDX-License-Identifier: BSD-3-Clause

set -euo pipefail

PLIST_NAME="io.cylonix.sase.direct.daemon"
PLIST_DST="/Library/LaunchDaemons/$PLIST_NAME.plist"
NOTIFIER_PLIST_NAME="io.cylonix.sase.direct.notifier"
NOTIFIER_PLIST_DST="/Library/LaunchAgents/$NOTIFIER_PLIST_NAME.plist"
APP_PATH="/Applications/Cylonix.app"
DAEMON_PATH="$APP_PATH/Contents/Resources/cylonixd"
NOTIFIER_PATH="$APP_PATH/Contents/Resources/CylonixNotifier.app/Contents/MacOS/CylonixNotifier"

# Create required directories
mkdir -p /var/lib/cylonix
mkdir -p /var/run/cylonix
mkdir -p /var/log/cylonix
chmod 755 /var/log/cylonix
# Ensure log file is world-readable so the app can display daemon logs
touch /var/log/cylonix/cylonixd.log
chmod 644 /var/log/cylonix/cylonixd.log

# Stop any existing daemon and WAIT for launchd to fully drop the label before we
# rewrite the plist. Deregister BY LABEL (bootout) rather than by plist path
# (`launchctl unload "$PLIST_DST"`): unload-by-path silently fails when a prior
# uninstall already removed the plist, leaving the label registered. bootout is
# asynchronous for a KeepAlive job, so poll until the label is actually gone —
# otherwise the bootstrap below hits EIO ("already loaded") and the daemon never
# starts (the uninstall→reinstall failure mode).
launchctl bootout "system/$PLIST_NAME" 2>/dev/null || true
i=0; while launchctl print "system/$PLIST_NAME" >/dev/null 2>&1 && [ "$i" -lt 100 ]; do
  i=$((i + 1)); sleep 0.1
done

# Install launchd plist
cat > "$PLIST_DST" <<'PLIST'
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>io.cylonix.sase.direct.daemon</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Applications/Cylonix.app/Contents/Resources/cylonixd</string>
		<string>--tun</string>
		<string>utun</string>
		<string>--state</string>
		<string>/var/lib/cylonix/cylonixd.state</string>
		<string>--statedir</string>
		<string>/var/lib/cylonix</string>
		<string>--socket</string>
		<string>/var/run/cylonix/cylonixd.sock</string>
		<string>--port</string>
		<string>41741</string>
	</array>
	<key>EnvironmentVariables</key>
	<dict>
		<key>CYLONIX_TAILDROP_DEFAULT_DOWNLOADS</key>
		<string>1</string>
	</dict>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>StandardOutPath</key>
	<string>/var/log/cylonix/cylonixd.log</string>
	<key>StandardErrorPath</key>
	<string>/var/log/cylonix/cylonixd.log</string>
</dict>
</plist>
PLIST

chown root:wheel "$PLIST_DST"
chmod 644 "$PLIST_DST"

# Load and start the daemon by LABEL. Do NOT suppress errors — a failure must be
# visible in /var/log/install.log, not a silently dead daemon. If the first
# bootstrap fails (e.g. a stubborn stale registration surviving the stop above),
# bootout + wait once more and retry before falling back to the legacy load.
if ! launchctl bootstrap system "$PLIST_DST"; then
  echo "postinstall: bootstrap failed; clearing stale/ghost registration and retrying" >&2
  # EIO with `print` reporting the label absent is a launchd "ghost": bootout has
  # nothing to grab but bootstrap collides with a cached entry. Evict it via every
  # available path — bootout by label, bootout by plist path, and the legacy
  # `remove` (a different code path that can drop a cached job bootout can't) —
  # then clear any disable override before retrying.
  launchctl bootout "system/$PLIST_NAME" 2>/dev/null || true
  launchctl bootout system "$PLIST_DST" 2>/dev/null || true
  launchctl remove "$PLIST_NAME" 2>/dev/null || true
  launchctl enable "system/$PLIST_NAME" 2>/dev/null || true
  i=0; while launchctl print "system/$PLIST_NAME" >/dev/null 2>&1 && [ "$i" -lt 100 ]; do
    i=$((i + 1)); sleep 0.1
  done
  launchctl bootstrap system "$PLIST_DST" || launchctl load "$PLIST_DST" || true
fi
launchctl kickstart -kp "system/$PLIST_NAME" 2>/dev/null || true

# Symlink CLI to /usr/local/bin so users can run "cylonix" from the terminal
mkdir -p /usr/local/bin
ln -sf "$APP_PATH/Contents/Resources/cylonix" /usr/local/bin/cylonix

# Note: the daemon resolves the active GUI/console user itself for
# Taildrop's Downloads/Cylonix destination (see
# tailscale/feature/taildrop/paths_cylonix.go), gated on
# CYLONIX_TAILDROP_DEFAULT_DOWNLOADS=1 (set above). No postinstall
# `cylonix set --operator` is required; admins on a shared Mac can still
# pin a specific user explicitly via `cylonix set --operator <user>`.

# Install the per-user notifier LaunchAgent. The system-wide LaunchAgents
# directory loads automatically for every user at login; we additionally
# bootstrap it into the currently-logged-in GUI session so notifications
# start working immediately without requiring logout/login.
cat > "$NOTIFIER_PLIST_DST" <<NOTIFIERPLIST
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>$NOTIFIER_PLIST_NAME</string>
	<key>ProgramArguments</key>
	<array>
		<string>$NOTIFIER_PATH</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>ThrottleInterval</key>
	<integer>5</integer>
	<key>StandardOutPath</key>
	<string>/tmp/cylonix-notifier.log</string>
	<key>StandardErrorPath</key>
	<string>/tmp/cylonix-notifier.log</string>
</dict>
</plist>
NOTIFIERPLIST
chown root:wheel "$NOTIFIER_PLIST_DST"
chmod 644 "$NOTIFIER_PLIST_DST"

GUI_USER="$(stat -f '%Su' /dev/console 2>/dev/null || true)"
if [[ -n "${GUI_USER:-}" && "$GUI_USER" != "root" && "$GUI_USER" != "_"* ]]; then
  GUI_UID="$(id -u "$GUI_USER" 2>/dev/null || true)"
  if [[ -n "${GUI_UID:-}" ]]; then
    # Replace any stale registration; ignore errors from never-loaded.
    launchctl bootout "gui/$GUI_UID/$NOTIFIER_PLIST_NAME" 2>/dev/null || true
    launchctl bootstrap "gui/$GUI_UID" "$NOTIFIER_PLIST_DST" 2>/dev/null || true
    launchctl enable "gui/$GUI_UID/$NOTIFIER_PLIST_NAME" 2>/dev/null || true
  fi
fi

# Kill any running Cylonix app so the new version is picked up
killall Cylonix 2>/dev/null || true
killall "cylonix-direct" 2>/dev/null || true

# Register app with LaunchServices so macOS uses this copy for share extension
LSREGISTER="/System/Library/Frameworks/CoreServices.framework/Versions/A/Frameworks/LaunchServices.framework/Versions/A/Support/lsregister"

# Unregister any stale copies from derived-data directories
"$LSREGISTER" -u "/tmp/cylonix-direct-derived" 2>/dev/null || true

# Register the installed app
"$LSREGISTER" -f "$APP_PATH"

# Register share extension with pluginkit
SHARE_EXT="$APP_PATH/Contents/PlugIns/ShareExtension.appex"
if [[ -d "$SHARE_EXT" ]]; then
  pluginkit -a "$SHARE_EXT" 2>/dev/null || true
fi

exit 0
